# snapraidd.service
[Unit]
Description=SnapRAID Daemon

# Ensure the network is up (for API/Notifications) and disks are mounted
After=network.target local-fs.target

# If SnapRAID requires specific mounts, you can list them here:
# RequiresMountsFor=/mnt/data1 /mnt/parity1

[Service]
Type=exec
# Path to your compiled binary
ExecStart=/usr/bin/snapraidd --config /etc/snapraidd.conf

# This implements the 'systemctl reload' command using SIGHUP
ExecReload=/bin/kill -HUP $MAINPID

# Restart logic: if it crashes, wait 5 seconds and try again
Restart=on-failure
RestartSec=5s

# Prevent the daemon from seeing /home to protect user privacy
ProtectHome=true

# Make /usr, /boot, etc. read-only for the daemon
ProtectSystem=full

# Give it its own private /tmp space
PrivateTmp=true

[Install]
# Start the daemon when the system reaches normal multi-user mode
WantedBy=multi-user.target
